# dockerhub.timer
[Unit]
Description=Runs docker updater everyday

[Timer]
OnCalendar=*-*-* 01:00:00
Unit=dockerhub.service

[Install]
WantedBy=multi-user.target
